# crontab/cron
* * * * * www-data /usr/local/bin/lockrun --lockfile=/home/www/metricsmaster.pl/spider.lockrun --quiet -- /home/www/metricsmaster.pl/html/spider.php
